# docker/github-runner/services.d/docker/run
#!/usr/bin/env bash
exec s6-notifyoncheck dockerd-rootless.sh --experimental --storage-driver fuse-overlayfs
