# docker/github-runner/services.d/runner/run
#!/usr/bin/env bash
s6-svwait -U /etc/services.d/docker
exec run.sh
